# deploy/systemd/jenn-mesh-dashboard.service
[Unit]
Description=JennMesh Dashboard (FastAPI)
Documentation=https://dev.azure.com/magnivation/JennMesh
After=network-online.target jenn-mesh-broker.service
Wants=network-online.target
Requires=jenn-mesh-broker.service

[Service]
Type=simple
User=jenn-mesh
Group=jenn-mesh
EnvironmentFile=/etc/jenn-mesh/env
ExecStart=/opt/jenn-mesh/current/venv/bin/uvicorn \
    jenn_mesh.dashboard.app:create_app \
    --factory \
    --host 0.0.0.0 \
    --port 8002 \
    --workers 2
WorkingDirectory=/opt/jenn-mesh/current
Restart=on-failure
RestartSec=5

# Resource limits — FastAPI with 2 workers
MemoryMax=256M
CPUQuota=50%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/jenn-mesh /var/log/jenn-mesh
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jenn-mesh-dashboard

[Install]
WantedBy=multi-user.target
